# MountFS/mount-internals.conf
description   "Mount Internal Partitions on ChromeOS"
author        "developer.kishor@gmail.com"

start on start-user-session

script

        exec /bin/mount-internals

end script
